# portail_captif/portail_captif.socket
[Unit]
Description=Socket Portail Captif

[Socket]
ListenStream=/tmp/gunicorn-portail_captif.sock

[Install]
WantedBy=sockets.target
